# sparql-server.service
[Unit]
Description=MCP SPARQL Server
After=network.target

[Service]
Type=simple
User=mcp
ExecStart=/usr/bin/python3 -m mcp_server_sparql.server --endpoint $SPARQL_ENDPOINT
Restart=on-failure
RestartSec=5
Environment=SPARQL_ENDPOINT=https://data.legilux.public.lu/sparqlendpoint
# Additional environment variables can be set here
#Environment=SPARQL_CACHE_ENABLED=true
#Environment=SPARQL_CACHE_TTL=300
#Environment=SPARQL_FORMAT=simplified

# Set the working directory if needed
#WorkingDirectory=/opt/mcp

# Make sure the process is always restarted
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target